# Nokia G-010S-A Official Firmware/squashfs-root/etc/init.d/diaglog.sh
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

#RCR ALU02298390 Resource Monitor
STOP=88
stop() {
	sync
	tar -zcf /logs/sysinfo.tgz -C /tmp sysinfo.log
}
